# deploy/crypto-signals-bot.service
[Unit]
Description=Crypto Signal Scanner Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/eva/workspace/crypto_agent
ExecStart=/usr/bin/python3 /home/eva/workspace/crypto_agent/main.py 15
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crypto-signals-bot

[Install]
WantedBy=default.target
